# infra/trustkit.service
[Unit]
Description=TrustKit (NDAI) API Server
After=network.target postgresql.service docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/trustkit
Environment=PATH=/opt/trustkit/.venv/bin:/usr/local/bin:/usr/bin
EnvironmentFile=/opt/trustkit/.env
ExecStart=/opt/trustkit/.venv/bin/uvicorn ndai.api.app:create_app --factory --host 127.0.0.1 --port 8100 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
